# systemd/linuxsecaudit.timer
[Timer]
OnBootSec=60
OnUnitActiveSec=3600
AccuracySec=600

[Install]
WantedBy=multi-user.target
