# systemd/linuxsecaudit.service
[Unit]
Description=Security auditing tool for Linux systems

[Service]
Type=oneshot
ExecStart=/usr/bin/linuxsecaudit
